# fs.conf
start on syno.network.ready
stop on runlevel [06]

respawn
respawn limit 5 10

setuid http

exec /usr/local/bin/fs 2332 /var/services/web
